# build/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:redis]
command=/etc/init.d/redis-server start

[program:mongodb]
command=/etc/init.d/mongod

[program:joola]
directory=/opt/joola
command=node joola.js
process_name=joola
user=joola
redirect_stderr=true