# build/docker/supervisord.conf
[supervisord]
nodaemon=true
redirect_stderr=true
stdout_logfile=/opt/joola/logs/supervisord.log

[program:sshd]
command=/usr/sbin/sshd -D
redirect_stderr=true
stdout_logfile=/opt/joola/logs/sshd.log

#[program:redis]
#command=/usr/bin/redis-server --bind 127.0.0.1 --dir /opt/joola/data/redis
#redirect_stderr=true
#stdout_logfile=/opt/joola/logs/redis.log

[program:mongodb]
command=/usr/bin/mongod --dbpath /opt/joola/data/mongodb
redirect_stderr=true
stdout_logfile=/opt/joola/logs/mongodb.log

#[program:rabbitmq]
#command=/usr/sbin/rabbitmq-server
#redirect_stderr=true
#stdout_logfile=/opt/joola/logs/rbq.log

[program:joola]
directory=/opt/joola/bin
command=sh run_within_docker.sh
startsecs=3
startretries=15
process_name=joola
user=joola
redirect_stderr=true
stdout_logfile=/opt/joola/logs/joola.log
